# etc/systemd/system/pikaur.service
[Unit]
Description=Update everything

[Service]
Type=oneshot
ExecStart=/usr/bin/pikaur -Syyu
